# automation/cursor-intelligence.service
[Unit]
Description=Cursor Intelligence Continuous Runner
After=network.target
Wants=network.target

[Service]
Type=simple
User=miami2
WorkingDirectory=/Users/miami2/ ZION_APP_AUG_11_25/zion.app
ExecStart=/usr/bin/node automation/cursor-intelligence-continuous-runner.cjs start
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cursor-intelligence

# Environment variables
Environment=NODE_ENV=production
Environment=PATH=/usr/local/bin:/usr/bin:/bin

# Resource limits
LimitNOFILE=65536
MemoryMax=512M

[Install]
WantedBy=multi-user.target
